# roles/install-etcd/files/etcd.service
[Unit]
Description=Etcd Server

[Service]
ExecStart=/usr/local/bin/etcd

[Install]
WantedBy=multi-user.target
